# eosplatform-gdrive-sync.service
[Unit]
Description=EOS Google Drive Sync Service
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/eosplatform
Environment="PATH=/usr/bin"
ExecStart=/home/ubuntu/eosplatform/eos_sync.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
